# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/supervisor

[program:llama-server]
command=/opt/llama/llama-server
    --model /models/Qwen3-VL-8B-Heretic-Stable.Q4_K_M.gguf
    --mmproj /models/Qwen3-VL-8B-Heretic-Stable.mmproj-bf16.gguf
    --host 0.0.0.0
    --port 10000
    --ctx-size 32768
    --n-gpu-layers 99
    --flash-attn on
    --parallel 1
    --jinja
directory=/opt/llama
autostart=true
autorestart=true
startsecs=10
startretries=20
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/llama-server.log
stderr_logfile=/var/log/supervisor/llama-server-error.log
environment=CUDA_VISIBLE_DEVICES="0"
priority=10

[program:open-webui]
command=/bin/bash /app/backend/start.sh
directory=/app/backend
autostart=true
autorestart=true
startsecs=10
startretries=20
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/open-webui.log
stderr_logfile=/var/log/supervisor/open-webui-error.log
environment=PORT="3000",HOST="0.0.0.0",DATA_DIR="/app/backend/data",ENABLE_OLLAMA_API="false",OPENAI_API_BASE_URL="http://127.0.0.1:10000/v1",OPENAI_API_KEY="gpu-game-local",WEBUI_AUTH="true"
priority=20
