# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:llama-server]
command=%(ENV_LLAMA_COMMAND)s
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/var/log/llama-server.log
stderr_logfile=/var/log/llama-server.err.log
redirect_stderr=false

[program:open-webui]
command=/opt/open-webui/venv/bin/open-webui serve --host 0.0.0.0 --port 3000
directory=/opt/open-webui
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/var/log/open-webui.log
stderr_logfile=/var/log/open-webui.err.log
redirect_stderr=false
environment=WEBUI_SECRET_KEY_FILE="/opt/gpu-game/.webui_secret_key",WEBUI_AUTH="True",OPENAI_API_BASE_URL="http://127.0.0.1:10000/v1",OPENAI_API_KEY="sk-local"
