# supervisord.conf
[unix_http_server]
file=/var/run/gpu-game-supervisor.sock
chmod=0700

[supervisord]
nodaemon=false
user=root
logfile=/var/log/gpu-game/supervisord.log
logfile_maxbytes=20MB
logfile_backups=3
pidfile=/var/run/gpu-game-supervisord.pid
childlogdir=/var/log/gpu-game

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/gpu-game-supervisor.sock


; ============================================================
; Основная диалоговая модель
; ============================================================

[program:rocinante]
command=/bin/bash -lc 'if [ "%(ENV_START_ROCINANTE)s" = "true" ]; then exec /app/llama-server --model "${ROCINANTE_FILE}" --alias Rocinante --host 127.0.0.1 --port 10000 --ctx-size 32768 --n-gpu-layers 99 --flash-attn on --parallel 1 --jinja; else exec sleep infinity; fi'
autostart=true
autorestart=true
startsecs=10
startretries=5
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/gpu-game/rocinante.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
redirect_stderr=true


; ============================================================
; Модель анализа изображений
; ============================================================

[program:vision]
command=/bin/bash -lc 'if [ "%(ENV_START_VISION)s" = "true" ]; then exec /app/llama-server --model "${VISION_FILE}" --mmproj "${VISION_MMPROJ_FILE}" --alias Qwen3-VL --host 127.0.0.1 --port 10001 --ctx-size 8192 --n-gpu-layers 99 --flash-attn on --parallel 1 --jinja; else exec sleep infinity; fi'
autostart=true
autorestart=true
startsecs=10
startretries=5
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/gpu-game/vision.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
redirect_stderr=true


; ============================================================
; Open WebUI
; ============================================================

[program:open-webui]
command=/bin/bash -lc 'if [ "%(ENV_START_WEBUI)s" = "true" ]; then exec /opt/open-webui/venv/bin/open-webui serve --host 127.0.0.1 --port 3000; else exec sleep infinity; fi'
directory=/opt/open-webui
autostart=true
autorestart=true
startsecs=10
startretries=5
stopasgroup=true
killasgroup=true
stopsignal=TERM
environment=
    DATA_DIR="/data/open-webui",
    WEBUI_AUTH="True",
    ENABLE_OLLAMA_API="False",
    WEBUI_SECRET_KEY_FILE="${WEBUI_SECRET_FILE}",
    OPENAI_API_BASE_URLS="http://127.0.0.1:10000/v1;http://127.0.0.1:10001/v1;http://127.0.0.1:10100/v1",
    OPENAI_API_KEYS="sk-local;sk-local;sk-local"
stdout_logfile=/var/log/gpu-game/open-webui.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
redirect_stderr=true


; ============================================================
; ComfyUI
; ============================================================

[program:comfyui]
command=/bin/bash -lc 'if [ "%(ENV_START_COMFY)s" = "true" ]; then cd /opt/ComfyUI && exec /opt/ComfyUI/venv/bin/python main.py --listen 127.0.0.1 --port 8188; else exec sleep infinity; fi'
directory=/opt/ComfyUI
autostart=true
autorestart=true
startsecs=10
startretries=5
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/gpu-game/comfyui.log
stdout_logfile_maxbytes=30MB
stdout_logfile_backups=3
redirect_stderr=true


; ============================================================
; Caddy: HTTPS и reverse proxy
; ============================================================

[program:caddy]
command=/bin/bash -lc 'if [ "%(ENV_START_CADDY)s" = "true" ]; then exec /usr/bin/caddy run --config "${ACTIVE_CADDYFILE}" --adapter caddyfile; else exec sleep infinity; fi'
autostart=true
autorestart=true
startsecs=5
startretries=5
stopasgroup=true
killasgroup=true
stopsignal=TERM
environment=
    XDG_DATA_HOME="/data/caddy",
    XDG_CONFIG_HOME="/config/caddy"
stdout_logfile=/var/log/gpu-game/caddy.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
redirect_stderr=true
